# init/base.crontab
#Telegram
#*/5 * * * 1-7  . /home/bot4voc/source/adminscript/crontabconf.conf ; $HOME/source/crontabconf/gitnodejs.sh t $LOGNAME >> /home/$LOGNAME/log/gitnodejs.log 2>&1
*/1 * * * 1-7  . /home/bot4voc/source/adminscript/crontabconf.conf ; $HOME/source/crontabconf/gitnodejs.sh t $LOGNAME >> /home/$LOGNAME/log/gitnodejs.log 2>&1


#Ngrok
#*/15 * * * 1-7 . $HOME/source/adminscript/bin/ngrok.sh 2>&1
*/1 * * * 1-7 . $HOME/source/adminscript/bin/ngrok.sh 2>&1
