# base.crontab
#Telegram
*/5 * * * 1-7  . $HOME/source/crontabconf/gitnodejs.sh t $LOGNAME >> /var/log/$LOGNAME/gitnodejs.log 2>&1


#Ngrok
*/15 * * * 1-7 . $HOME/source/crontabconf/ngrok.sh

